# docker/data/supervisord.conf
[supervisord]
nodaemon=true

[program:tightvnc]
command=su ubuntu -c "/usr/bin/tightvncserver :1" 
user=ubuntu
autorestart=true

[program:nginx]
command=/usr/sbin/nginx -g 'daemon off;'
autorestart=true

[program:wetty]
command=/usr/bin/node app.js -p 3000
directory=/opt/wetty
autorestart=true

[program:ipaddress]
command=/usr/local/bin/caddy
directory=/
autorestart=true
